# etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bash

dnsmasq --no-daemon --dhcp-range=${DHCP_RANGE_START},proxy ${DNSMASQ_ARGS}
